# systemd/ydotoold.service
[Unit]
Description=ydotool daemon (user)

[Service]
ExecStart=/usr/bin/ydotoold

[Install]
WantedBy=default.target
